# config/mjc-website.service
[Unit]
Description=MJC Website Production Server
After=network.target mysql.service
Requires=mysql.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/MJCWebsite
ExecStart=/path/to/MJCWebsite/scripts/start.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:/path/to/MJCWebsite/logs/systemd-output.log
StandardError=append:/path/to/MJCWebsite/logs/systemd-error.log

# Environment
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target
